# chef/cookbooks/git/templates/default/upstart.conf.erb
description "<%= @service_name %> server"
author "git-barclamp <root@crowbar>"

start on (filesystem and net-device-up IFACE!=lo)
stop on runlevel [016]

<% if @service_name.include? "nova" -%>
chdir /var/run

pre-start script
    mkdir -p /var/run/nova
    chown nova:root /var/run/nova/

    mkdir -p /var/lock/nova
    chown nova:root /var/lock/nova/
    <%= "modprobe nbd" if @service_name == "nova-compute" %>
end script
<% params = "--config-file=/etc/nova/nova.conf" -%>
<% params += " --config-file=/etc/nova/nova-compute.conf" if @service_name == "nova-compute" -%>
<% else -%>
<% params = nil -%>
respawn
<% end %>

<% params += " #{@opt_params}" if @opt_params -%>
<% cmd = @bin_name -%>
<% cmd = "#{@opt_path}/#{@bin_name}" if @opt_path -%>
<% cmd += " #{params}" if params -%>
exec su -s /bin/sh -c "exec <%= cmd %>" <%= @user %>
